# infra/quadlets/wayofmono-caddy.container
[Unit]
Description=WayOfMono Caddy Reverse Proxy
After=wayofmono-nextjs.service

[Container]
Image=caddy:2-alpine
ContainerName=wayofmono-caddy
PublishPort=8090:8090
Volume=<REPO_PATH>/infra/podman/Caddyfile:/etc/caddy/Caddyfile:ro,Z
Volume=wayofmono-caddy-data.volume:/data
Volume=wayofmono-caddy-config.volume:/config
Network=wayofmono.network

[Service]
Restart=always

[Install]
WantedBy=multi-user.target default.target
